# buildroot-external/board/raspberrypi/raspberrypi0w/inittab
## insert into /etc/inittab to load the wifi before the networkscripts
# Load Wifi driver
::sysinit:/sbin/modprobe brcmfmac

# now run any rc scripts
::sysinit:/etc/init.d/rcS


sed -i '/# now run any rc scripts/::sysinit:/sbin/modprobe brcmfmac' /etc/inittab

